# services/io.kud.ai-feed-digest.publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.kud.ai-feed-digest.publish</string>

  <key>UserName</key>
  <string>kud</string>

  <key>WorkingDirectory</key>
  <string>/Users/kud/Projects/ai-feed-digest</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/kud</string>
    <key>PATH</key>
    <string>/Users/kud/.npm-global/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>set -eu; set -o pipefail; STASH_REF=""; if ! /usr/local/bin/git diff --quiet || ! /usr/local/bin/git diff --quiet --cached; then /usr/local/bin/git stash push -u -m "ai-feed-digest.publish" >/dev/null; STASH_REF=$(/usr/local/bin/git stash list | head -n1 | cut -d: -f1); fi; cleanup() { if [ -n "$STASH_REF" ]; then /usr/local/bin/git stash pop --index "$STASH_REF" >/dev/null 2>&1 || /usr/local/bin/git stash pop "$STASH_REF" >/dev/null 2>&1 || true; fi; }; trap cleanup EXIT; /usr/local/bin/git checkout generated; /usr/local/bin/npm run publish; /usr/local/bin/git checkout -</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>35</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/kud/Library/Logs/ai-feed-digest.publish.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kud/Library/Logs/ai-feed-digest.publish.error.log</string>
</dict>
</plist>
